# supervisord.conf
[unix_http_server]
file=./tmp/supervisor.sock                       ; path to your socket file
# chmod=0777                 ; socket file mode (default 0700)
# sockchown=liuzhaoming:root       ; socket file uid:gid owner
username=user              ; (default is no username (open server))
password=123               ; (default is no password (open server))

[inet_http_server]
# Web管理界面设定
port=127.0.0.1:9001 ;
username = user ;
password = 123 ;

[supervisord]
logfile=./logs/supervisord.log                    ; supervisord log file
logfile_maxbytes=10MB                           ; maximum size of logfile before rotation
logfile_backups=5                              ; number of backed up logfiles
loglevel=info                                   ; info, debug, warn, trace
pidfile=./tmp/supervisord.pid                    ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=liuzhaoming                                   ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://./tmp/supervisor.sock           ; use a unix:// URL  for a unix socket
username=user
password=123

[program:appium_server]
command=/Users/liuzhaoming/software/python/.virtualenvs/appium_server/bin/python  ./appium_server.py --production=false ;
autostart=true ;
process_name=%(program_name)s  ;
redirect_stderr=true  ;
stopsignal=QUIT  ;
stdout_logfile= ./logs/supervisord_server.log  ;